# systemd/honeypy.service
[Unit]
Description=HONEPY.

[Service]
ExecStart=/usr/bin/python3 /main.py
Restart=always

[Install]
WantedBy=multi-user.target